# systemd/dexter-rpi-bridge.service
[Unit]
Description=Dexter RPi Hardware Bridge (roslibpy → PCA9685 → servos)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi3
WorkingDirectory=/home/pi3
EnvironmentFile=/home/pi3/dexter_bridge.env
ExecStart=/usr/bin/python3 /home/pi3/dexter_rpi_bridge/dexter_rpi_bridge/hardware_node.py \
          --host ${LAPTOP_IP} --port ${ROSBRIDGE_PORT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
